# opencode/.config/systemd/user/issue-resolver.service
[Unit]
Description=Issue Resolver — auto-resolve @yesman GitHub issues via AI
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/share/mise/shims/bun run %h/.config/opencode/scripts/issue-resolver.ts
WorkingDirectory=%h
Environment=PATH=%h/.local/share/mise/shims:%h/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
