# opencode/.config/systemd/user/obsidian-vault-sync.service
[Unit]
Description=Obsidian Vault Git Sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/obsidian-vault-sync.sh
Environment=HOME=%h

[Install]
WantedBy=default.target
