# opencode/.config/systemd/user/obsidian-vault-sync.timer
[Unit]
Description=Run Obsidian Vault Git Sync every 15 minutes

[Timer]
OnCalendar=*:0/15
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
