# opencode/.config/systemd/user/obsidian-vault-sync.service
[Unit]
Description=Obsidian Vault continuous sync
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Documents/MyObsidianVault
ExecStart=/usr/bin/timeout --kill-after=30s 6h %h/.bun/bin/ob sync --path %h/Documents/MyObsidianVault --continuous
Restart=always
RestartSec=10
KillMode=mixed
TimeoutStopSec=45
Environment=HOME=%h
Environment=PATH=%h/.local/share/mise/shims:%h/.bun/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
