# opencode/.config/systemd/user/pr-dashboard.timer
[Unit]
Description=Run PR Dashboard every morning

[Timer]
OnCalendar=*-*-* 08:30:00
Persistent=true

[Install]
WantedBy=timers.target
