# opencode/.config/systemd/user/pr-dashboard.timer
[Unit]
Description=Run PR Dashboard on weekday mornings

[Timer]
OnCalendar=Mon..Fri *-*-* 08:30:00
Persistent=true

[Install]
WantedBy=timers.target
