# opencode/.config/systemd/user/yesmand.service
[Unit]
Description=Yesman Automation Daemon + UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=PATH=%h/.bun/bin:%h/.cargo/bin:%h/.local/share/mise/shims:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=YESMAN_UI_ADDR=0.0.0.0:3777
EnvironmentFile=%h/.config/yesman.env
WorkingDirectory=%h
ExecStartPre=/bin/sh -lc 'case "$$YESMAN_UI_PASSWORD" in ""|changeme) echo "Set YESMAN_UI_PASSWORD in %h/.config/yesman.env before starting yesmand.service" >&2; exit 1;; esac'
ExecStart=/bin/sh -lc 'mkdir -p %h/.local/share/yesman/runtime && cd %h/.local/share/yesman/runtime && exec %h/.cargo/bin/yesman up'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
