# yesmand/ops/systemd/yesmand-pr-reviews.service
[Unit]
Description=Yesman PR Reviews Automation
After=network-online.target opencode-automation.service
Wants=network-online.target opencode-automation.service

[Service]
Type=oneshot
WorkingDirectory=%h/dotfiles/yesmand
Environment=PATH=%h/.local/share/mise/shims:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=%h/.config/opencode/.env
ExecStartPre=/bin/sh -c 'for i in $(seq 1 30); do code=$(curl -s -o /dev/null -w "%%{http_code}" http://127.0.0.1:4097/global/health || true); if [ "$code" = "200" ] || [ "$code" = "401" ]; then exit 0; fi; sleep 0.5; done; echo "opencode-automation on 4097 not ready" >&2; exit 1'
ExecStart=%h/.local/share/mise/installs/bun/1.3.9/bin/bun run src/index.ts --once --config %h/dotfiles/yesmand/config.pr-reviews.json

[Install]
WantedBy=default.target
