# yesmand/ops/systemd/yesmand-pr-reviews.timer
[Unit]
Description=Run yesmand PR reviews hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
